# deploy/backplane.service.tmpl
[Unit]
Description=Backplane MCP Server
After=network.target

[Service]
Type=simple
User=%%SERVICE_USER%%
WorkingDirectory=%%INSTALL_DIR%%
ExecStart=%%INSTALL_DIR%%/.venv/bin/python -m backplane.mcp
Restart=on-failure
RestartSec=5
EnvironmentFile=%%INSTALL_DIR%%/.env
StandardOutput=append:%%LOG_DIR%%/backplane.log
StandardError=append:%%LOG_DIR%%/backplane.log

[Install]
WantedBy=multi-user.target
